# deployment/production/episciences-next.service
[Unit]
Description=Episciences Next.js Frontend
After=network.target

[Service]
Type=simple
User=epi-preprod-deploy
Group=www-data
# ansistrano deploys to current/ symlink
WorkingDirectory=/var/www/episciences-front-next/current
ExecStart=/usr/bin/node /var/www/episciences-front-next/current/.next/standalone/server.js
Restart=always
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
# Optional: Override API target if needed
# Environment=API_PROXY_TARGET=https://api.episciences.org/api

[Install]
WantedBy=multi-user.target